# config/supervisord/rest2web.conf
[program:rest2web]
command = /usr/local/bin/rest2web
user = xavinsky
stderr_logfile = /var/log/rest2web.err
stdout_logfile = /var/log/rest2web.log
directory = /var/www/rest2web
