# config/supervisord/webinrest.conf
[program:webinrest]
command = /usr/local/bin/webinrest
user = xavinsky
stderr_logfile = /var/log/webinrest.err
stdout_logfile = /var/log/webinrest.log
directory = /var/www/webinrest
